# config/systemd/hkn-rails-logrotate.timer
[Unit]
Description=Run hkn-rails logrotate daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
